# ansible/roles/dyndns/files/dns-update.cron
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin:/usr/local/bin
*/3 * * * * root /opt/commapps/dyndns/dns-update.sh
